# templates/supervisord.conf
[program:sshtron]
command=/opt/sshtron/bin/sshtron
directory=/opt/sshtron/bin
environment=SSH_PORT="{{ sshtron_ssh_port }}",PORT="{{ sshtron_http_port }}"